# phantomjs/files/etc/supervisord.d/tcpdump.ini
; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:tcpdump]
command=/usr/sbin/tcpdump -nn ; the program (relative uses PATH, can take args)
autostart=true                ; start at supervisord start (default: true)
autorestart=false              ; retstart at unexpected quit (default: true)
priority=250
exitcodes=0,2
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/dev/fd/1        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=0   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile=/dev/fd/1        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=0   ; max # logfile bytes b4 rotation (default 50MB)
